# scheduling/cron.d/crontab
* * * * * /data/tasks/sync-contact --target ${MARIADB_DSN} --source ${POSTGRESQL_DSN}
