# scheduling/cron.d/crontab
* * * * * /data/tasks/sync-contact partial
