# scheduling/cron.d/crontab
* * * * * /data/tasks/sync-contact --target mariadb://user:password@host:port/dbname --source postgres://user:password@host:port/dbname
